# decompress-kernel.service
[Unit]
Description="Decompress vmlinuz into bootable vmlinux file"

[Service]
ExecStart=/usr/local/bin/decompress-kernel